# .devcontainer/cbpi-default-dev-config/craftbeerpi.service
[Unit]
Description=Craftbeer Pi

[Service]
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/cbpi start

[Install]
WantedBy=multi-user.target
